# etc/profile.d/01-cargo.rc
if(~ $path *^$home/.cargo/bin^*) {}
if not {
	path=($path $home/.cargo/bin)
}
