# etc/profile.d/10-editor.rc
if(~ $#DISPLAY 0) {
	EDITOR=(sam -d)
}
if not {
	EDITOR=E
}
