# etc/profile.d/10-logname.sh
if logname > /dev/null 2>&1
then
	LOGNAME=$(logname)
else
	LOGNAME=$(basename $HOME)
fi
export LOGNAME
