# etc/profile.d/10-logname.rc
if(logname >/dev/null >[2=1]) {
	LOGNAME=`{logname}
}
if not {
	LOGNAME=`{basename $home}
}
