# etc/profile.d/10-xdg.rc
XDG_CONFIG_HOME=$home/etc
XDG_CACHE_HOME=$home/var/cache
XDG_DATA_HOME=$home/var
