# etc/profile.d/99-my.rc
if(~ $path *^$home/bin^*) {}
if not {
	path=($home/bin $path)
}
