# etc/profile.d/env.sh
if [ -r "$HOME/etc/shrc" ]
then
	ENV="$HOME/etc/shrc"
	export ENV
fi
