# etc/profile.d/go.sh
if type go > /dev/null 2>&1
then
	GOPATH="$HOME"
	export GOPATH
fi
