# etc/profile.d/visual.sh
if type vis >/dev/null 2>&1
then
	VISUAL=vis
else
	VISUAL=vi
fi
export VISUAL
