# etc/profile.d/xdg.sh
XDG_CONFIG_HOME=$HOME/etc
XDG_CACHE_HOME=$HOME/var/cache
XDG_DATA_HOME=$HOME/lib/share
export XDG_CONFIG_HOME XDG_CACHE_HOME XDG_DATA_HOME
